# files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

nvram set partialboots=0 && nvram commit
/etc/fix_interfaces.sh

exit 0

